# meta-my-iotc-python-sdk-example/recipes-systemd/iotc-telemetry-demo-service/files/iotc-telemetry-demo.service
[Unit]
Description= iotconnect python service, runs forever
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 -u /usr/bin/iotc/telemetry_demo.py /usr/bin/iotc/credentials.json
Restart=always
StandardOutput=journal+console
StandardError=journal+console

[Install]
WantedBy=multi-user.target